# runtime/launchd/com.memory-index-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.local.memory-index-watcher</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>NVM="$HOME/.nvm/versions/node"; L=$(ls -1 "$NVM" 2>/dev/null | sort -V | tail -1); [ -n "$L" ] &amp;&amp; export PATH="$NVM/$L/bin:$PATH"; [ -d /opt/homebrew/bin ] &amp;&amp; export PATH="/opt/homebrew/bin:$PATH"; exec node "$HOME/.bin/memory-index-build.mjs" --watch</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>$HOME/.local/share/journal/.memory-index-watch.log</string>
  <key>StandardErrorPath</key><string>$HOME/.local/share/journal/.memory-index-watch.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
